# roles/display/templates/geoclue-agent.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run Geoclue agent required by applications

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
