# roles/ledger/templates/ledger-backup.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Ledger backup

[Service]
ExecStart={{ scripts_directory }}/bledger
Type=oneshot
