# roles/ledger/templates/ledger.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Add recurring transactions to ledger

[Timer]
OnCalendar=*-*-* 23:00:00

[Install]
WantedBy=timers.target
